# systemd/one-million-sliders.socket
[Unit]

[Socket]
ListenStream=127.0.0.1:8081
ReusePort=true

[Install]
WantedBy=sockets.target
